# com.gentlyyonder.gsc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.gentlyyonder.gsc</string>
  <key>ProgramArguments</key><array><string>/Users/ky/travel-now-agent/run_gsc.sh</string></array>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>10</integer></dict>
  <key>StandardOutPath</key><string>/Users/ky/travel-now-agent/gsc_reports/launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/ky/travel-now-agent/gsc_reports/launchd.log</string>
</dict></plist>
